# database/backup/systemd/datenkraken_db_backup.service
[Unit]
description = "Backup db script"

[Service]
ExecStart=/home/$USER/backup.sh
